# pia-qbittorrent-sync.service
[Unit]
Description=PIA qBittorrent Sync Service
After=network-online.target wg-quick@wg0.service
Wants=network-online.target

[Service]
Type=simple
User=pia-sync
Group=pia-sync
WorkingDirectory=/opt/pia-qbittorrent-sync
ExecStart=/opt/pia-qbittorrent-sync/venv/bin/python3 /opt/pia-qbittorrent-sync/pia_qbittorrent_sync.py
Restart=always
RestartSec=10
# Exit code 78 (EX_CONFIG) indicates a configuration error (e.g. missing/insecure credentials).
# Prevent systemd from restarting in a tight loop on misconfiguration.
RestartPreventExitStatus=78

# Environment variables (customize these)
Environment="PIA_GATEWAY=10.0.0.1"
# PIA_CA_CERT is required. The service will not start without it.
# Download ca.rsa.4096.crt from the PIA manual connection guide.
Environment="PIA_CA_CERT=/etc/pia/ca.rsa.4096.crt"
# PIA_HOSTNAME should match the server name in the PIA CA certificate.
# Environment="PIA_HOSTNAME=your_server_hostname"
Environment="PIA_TOKEN_FILE=/run/pia/token"
Environment="QBITTORRENT_HOST=http://localhost:8080"
# WARNING: Replace CHANGE_ME with your actual qBittorrent credentials before starting the service.
# The service will refuse to start if these are left as CHANGE_ME or any other known default.
Environment="QBITTORRENT_USERNAME=CHANGE_ME"
Environment="QBITTORRENT_PASSWORD=CHANGE_ME"
Environment="CHECK_INTERVAL=300"
Environment="LOG_LEVEL=INFO"
Environment="LOG_FILE=/var/log/pia-qbittorrent-sync/pia_qbittorrent_sync.log"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/pia-qbittorrent-sync /run/pia
LogsDirectory=pia-qbittorrent-sync
LogsDirectoryMode=0750
RuntimeDirectory=pia
RuntimeDirectoryMode=0700

[Install]
WantedBy=multi-user.target
